# scripts/splunk
#!/bin/bash
#
#------------------------------------------------------
#   scriptname: /etc/init.d/splunk
#   $VERSION: [v1.0-21] $
#   $DATE:    [Wed Dec 27,2017 - 11:11:38AM -0600] $
#	$AUTHOR:  [mhassan2 <mhassan@splunk.com>] $
#
#------------------------------------------------------
#
# init script for Splunk.
# generated by 'splunk enable boot-start'.
#
### BEGIN INIT INFO
# Provides: splunkd
# Required-Start: $remote_fs
# Required-Stop: $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start splunk
# Description: Splunk indexer service
### END INIT INFO
#
SPLUNK_USER=splunk
SPLUNK_HOME=/opt/splunk
RETVAL=0

splunk_start() {
  #echo Starting Splunk...
  #"/opt/splunk/bin/splunk" start --no-prompt --answer-yes
  /sbin/splunk_first_run.sh start-service
  	printf "\033[1;33mSync with supervisord [supervisorctl start splunkd]: "
	supervisorctl start splunkd
	printf "\033[0m\n"

  RETVAL=$?
}
splunk_stop() {
  echo Stopping Splunk...
   	sudo -HEu ${SPLUNK_USER} sh -c "${SPLUNK_HOME}/bin/splunk stop"
  	printf "\033[1;33mSync with supervisord [supervisorctl stop splunkd]: "
	supervisorctl stop splunkd
	printf "\033[0m\n"
  RETVAL=$?
}
splunk_restart() {
  echo Restarting Splunk...
  splunk_stop
 # "/opt/splunk/bin/splunk" restart
#	sudo -HEu ${SPLUNK_USER} sh -c "${SPLUNK_HOME}/bin/splunk restart"
	splunk_start
  RETVAL=$?
}
splunk_status() {
  echo Splunk status:
  "/opt/splunk/bin/splunk" status
  RETVAL=$?
}
case "$1" in
  start)
    splunk_start
    ;;
  stop)
    splunk_stop
    ;;
  restart)
    splunk_restart
    ;;
  status)
    splunk_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac
trap "/opt/splunk/bin/splunk stop" SIGINT SIGTERM
#wait
#Keep running if started by supervisord, otherwise terminate
MY_PID=$$
SUPERVISOR_PIDS=$(supervisorctl pid all)
#echo "my pid: [$MY_PID]  suprverisord pids:[$SUPERVISOR_PIDS]"
if [[ "$SUPERVISOR_PIDS" =~ "$MY_PID" ]]; then
	while true; do continue;done
fi
#trap "sudo -HEu ${SPLUNK_USER} ${SPLUNK_HOME}/bin/splunk stop" SIGINT SIGTERM
#trap "/opt/splunk/bin/splunk stop" SIGINT SIGTERM
exit $RETVAL


